# root/config/service/nginx/run
#!/bin/sh
nginx -g 'pid /tmp/nginx.pid; daemon off; error_log /dev/stdout info;'
